# sui-intent-agent.service
[Unit]
Description=Sui Intent Agent - Next.js App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/sui-intent-agent
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3
StandardOutput=append:/var/log/sui-intent-agent.log
StandardError=append:/var/log/sui-intent-agent-error.log

[Install]
WantedBy=multi-user.target
